# container_startup.service
[Unit]
Description=Startup script to keep things running

[Service]
ExecStart=/usr/local/bin/startup.sh

[Install]
WantedBy=multi-user.target